# relay-ssh-agent.service
[Install]
WantedBy=default.target

[Service]
ExecStart=~/.local/relay-ssh-agent.sh
