# relay-ssh-agent.service
[Install]
WantedBy=default.target

[Service]
ExecStart=$HOME/.local/relay-ssh-agent.sh
