# relay-ssh-agent.service
[Install]
WantedBy=default.target

[Service]
ExecStart=%h/.local/relay-ssh-agent.sh
